# bin/bixby-agent.initd
#!/bin/bash

### BEGIN INIT INFO
# Provides:          bixby-agent
# Required-Start:    $local_fs $remote_fs $syslog $named $network $time
# Required-Stop:     $local_fs $remote_fs $syslog $named $network
# Should-Start:
# Should-Stop:
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Start/Stop the bixby-agent daemon
### END INIT INFO

# processname: bixby-agent
# pidfile: /opt/bixby/var/bixby-agent.pid

NAME=bixby-agent
DESC=bixby-agent

# must set in order to work
BIXBY_HOME=/opt/bixby
BIXBY_BIN=$BIXBY_HOME/bin/bixby-agent
BIXBY_PID=$BIXBY_HOME/var/bixby-agent.pid

case $1 in
start)
        echo -n "Starting $DESC: "
        $BIXBY_BIN start
        if [[ $? -eq 0 ]]; then
          echo "done"
        fi
        ;;
stop)
        if [[ -f $BIXBY_PID ]]; then
          echo -n "Stopping $DESC: "
        fi
        $BIXBY_BIN stop
        if [[ $? -eq 0 ]]; then
          echo "done"
        fi
        ;;
restart)
        echo -n "Restarting $DESC: "
        $BIXBY_BIN restart
        if [[ $? -eq 0 ]]; then
          echo "done"
        fi
        ;;
status)
        $BIXBY_BIN status
        ;;
zap)
        echo -n "Zapping $DESC: "
        $BIXBY_BIN zap
        if [[ $? -eq 0 ]]; then
          echo "done"
        fi
        ;;
*)
  echo "Usage: $NAME {start|stop|restart|status|zap}" >&2
  exit 1
  ;;
esac

exit 0
